# .services/drowsiness.service
[Unit]
Description=Drowsiness Detection FastAPI Service
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/drowsiness-detection
ExecStart=/opt/drowsiness-detection/scripts/start.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target